# temboard/temboard-agent/9.6/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgres]
command=/docker-entrypoint.sh postgres
autorestart=true

[program:temboard-agent]
command=/usr/local/src/temboard-master/temboard-agent/temboard-agent
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
